# backend/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid

; ── Apache (servidor HTTP principal) ─────────────────────────────────────────
[program:apache]
command=apache2-foreground
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/apache.log
stderr_logfile=/var/www/storage/logs/apache_error.log
priority=1

; ── Horizon (Supervisa todos los trabajos asíncronos y colas Redis) ──
[program:horizon]
command=php /var/www/artisan horizon
autostart=true
autorestart=true
user=www-data
numprocs=1
stdout_logfile=/var/www/storage/logs/horizon.log
stderr_logfile=/var/www/storage/logs/horizon_error.log
priority=2

; ── Laravel Reverb (WebSocket server) ────────────────────────────────────────
[program:reverb]
command=php /var/www/artisan reverb:start --host=0.0.0.0 --port=8081 --no-interaction
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/www/storage/logs/reverb.log
stderr_logfile=/var/www/storage/logs/reverb_error.log
priority=3
